# config/files/usr/etc/systemd/system/var-mnt-appdata.automount
[Unit]
Description=Automount ono data

[Automount]
Where=/var/mnt/appdata

[Install]
WantedBy=multi-user.target
