# config/files/usr/etc/systemd/system/var-mnt-data.automount
[Unit]
Description=Automount ono data

[Automount]
Where=/var/mnt/data

[Install]
WantedBy=multi-user.target
