# config/files/usr/etc/systemd/system/var-mnt-backups.automount
[Unit]
Description=Automount ono backups

[Automount]
Where=/var/mnt/backups

[Install]
WantedBy=multi-user.target
