# config/files/usr/etc/systemd/system/var-mnt-photos.automount
[Unit]
Description=Automount ono photos

[Automount]
Where=/var/mnt/photos

[Install]
WantedBy=multi-user.target
