# config/files/usr/etc/systemd/system/var-mnt-docker.automount
[Unit]
Description=Automount ono docker

[Automount]
Where=/var/mnt/docker

[Install]
WantedBy=multi-user.target
