# config/files/usr/etc/systemd/system/var-mnt-docs.automount
[Unit]
Description=Automount ono docs

[Automount]
Where=/var/mnt/docs

[Install]
WantedBy=multi-user.target
