# config/files/usr/lib/systemd/user/dotfiles-update.timer
[Unit]
Description=Pull dotfiles container image

[Timer]
OnBootSec=1h
OnUnitInactiveSec=1d

[Install]
WantedBy=timers.target
