# config/files/usr/lib/systemd/system/var-mnt-isos.automount
[Unit]
Description=Automount ono isos

[Automount]
Where=/var/mnt/isos

[Install]
WantedBy=multi-user.target
